# com.cf.agent.aiv2.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cf.agent.aiv2.bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/anthonypium/Desktop/CF_AI_Agent/cf_agent_aiv2_bot_minimal_reverted.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/anthonypium/Desktop/CF_AI_Agent</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/anthonypium/clawd/logs/cf_agent_bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/anthonypium/clawd/logs/cf_agent_bot_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>